# supervisord.conf
[supervisord]
nodaemon=true

[program:security_scheduler]
command=/bin/sh -l -c 'cd /app; bundle exec rake versioneye:scheduler_security_prod'
autorestart=true
